# python-sidecar/bakedbot-sidecar.service
[Unit]
Description=Markitbot Python Sidecar - NotebookLLM MCP Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bakedbot-sidecar
Environment="PATH=/opt/bakedbot-sidecar/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DISPLAY=:99"
ExecStart=/usr/bin/xvfb-run -a /opt/bakedbot-sidecar/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bakedbot-sidecar

[Install]
WantedBy=multi-user.target
